# runit/sv/pulseaudio/finish
#!/bin/sh -u
exec 2>&1
exec pulseaudio --kill
