# runit/sv/pulseaudio/run
#!/bin/sh -u

# DISPLAY is undefined at login.
export DISPLAY=:0

exec 2>&1
exec pulseaudio -v
